# deploy/systemd/xblog-web.service
[Unit]
Description=xblog Next.js frontend
After=network.target xblog-api.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/xblog/frontend
Environment=NODE_ENV=production
EnvironmentFile=-/srv/xblog/frontend/.env.local
ExecStart=/srv/xblog/frontend/node_modules/.bin/next start -H 127.0.0.1 -p 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
