# backend/deploy/supervisor/backend.conf
; supervisor 进程配置（部署前请按实际环境修改 directory / user / 路径）
[program:backend]
directory=/srv/backend
command=/srv/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/backend/logs/uvicorn.log
environment=APP_ENV="production"
